# bin/crontab
#!/bin/sh
# Ce script permet le relancement de sdreams lors d'un plantage. Si c'est un die, il ne relance pas
#
# $Id: crontab,v 1.3 2005/04/20 14:30:38 bugs Exp $
#

if test ! "$1" ; then
  cd ..;
else
  if [ -d "$1" ]; then
    cd $1;
  else
    echo "Le répertoire $1 n'exite pas";
    exit 0;
  fi
fi

if [ -f sdreams.pid ] && ! ps -p `cat sdreams.pid` >/dev/null 2>&1
then
    if [ -f logs/err.log ]; then echo `date` - Redemarrage de SDreams par crontab >> logs/err.log; fi
    rm sdreams.pid
    sleep 5
    ./SDreams
fi
